# .config/systemd/user/restic-backup.service
[Unit]
Description=Restic backup service

[Service]
ExecStart=/home/cyryl/tools/backup.sh

[Install]
WantedBy=default.target
